# ubuntu-cpu/scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter]
command=/bin/bash /run/scripts/jupyter-notebook.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
